# templates/graylog2_server.conf.erb
[program:graylog-server]
command=java -jar graylog2-server.jar 
directory=<%= scope.lookupvar('graylog2::params::graylog_server_dir') %>
user=root
stopsignal=QUIT
